# deploy/systemd/sarl-platform-status.timer
[Unit]
Description=Daily SARL platform status digest to Telegram

[Timer]
OnCalendar=*-*-* 08:00:00
AccuracySec=1min
Persistent=true
Unit=sarl-platform-status.service

[Install]
WantedBy=timers.target
