# systemd/sarl-agent-ai-telegram-work-sync.service
[Unit]
Description=Synchronize durable Telegram missions into Hermes Kanban
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/root/CascadeProjects/SARL-agent-ai/scripts/run-telegram-work-sync.sh
User=root
Group=root
Nice=10
